# contrib/mac_osx/org.boxbackup.bbackupd.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.boxbackup.bbackupd</string>
	<key>RunAtLoad</key>
	<true/>
	<key>ProgramArguments</key>
	<array>
		<string>@prefix@/sbin/bbackupd</string>
		<string>-F</string>
		<string>@prefix@/etc/boxbackup/bbackupd.conf</string>
	</array>
        <key>LowPriorityIO</key>
        <true/>
        <key>Nice</key>
        <integer>1</integer>
</dict>
</plist>
